# proxmox/host/node-exporter/node-exporter.service
[Unit]
Description=Prometheus node_exporter (host hardware metrics)
After=network-online.target
Wants=network-online.target

[Service]
User=node_exporter
Group=node_exporter
ExecStart=/usr/local/bin/node_exporter \
    --collector.textfile.directory=/var/lib/node_exporter/textfile \
    --collector.systemd \
    --web.listen-address=:9100
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
